# timebomb-cleaner.service
[Unit]
Description=Timebomb Cleaner: Delete expired files

[Service]
Type=oneshot
ExecStart=%h/.timebomb/cleaner.sh
