# timebomb-linux/timebomb-cleaner.timer
[Unit]
Description=Run Timebomb Cleaner every hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
